# infra/launchd/com.claudeclaw.consolidator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.claudeclaw.consolidator</string>

  <!-- Runs every 30 minutes: extractor then consolidator. -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>cd /Users/jaywest/claudeclaw &amp;&amp; .venv/bin/python memory/extractor.py &amp;&amp; .venv/bin/python memory/consolidator.py</string>
  </array>

  <key>StartInterval</key><integer>1800</integer>

  <key>RunAtLoad</key><true/>

  <key>StandardOutPath</key><string>/Users/jaywest/Library/Logs/claudeclaw/consolidator.out</string>
  <key>StandardErrorPath</key><string>/Users/jaywest/Library/Logs/claudeclaw/consolidator.err</string>
</dict>
</plist>
