# infra/launchd/com.claudeclaw.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.claudeclaw.dashboard</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>--import=tsx</string>
    <string>/Users/jay/claudeclaw/apps/dashboard/src/server.ts</string>
  </array>

  <key>WorkingDirectory</key><string>/Users/jay/claudeclaw/apps/dashboard</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>

  <key>StandardOutPath</key><string>/Users/jay/Library/Logs/claudeclaw/dashboard.out</string>
  <key>StandardErrorPath</key><string>/Users/jay/Library/Logs/claudeclaw/dashboard.err</string>
</dict>
</plist>
