# infra/launchd/com.claudeclaw.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>           <string>com.claudeclaw.dashboard</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/jaywest/claudeclaw/bin/launchd-wrap.sh</string>
    <string>dashboard</string>
  </array>

  <key>WorkingDirectory</key><string>/Users/jaywest/claudeclaw</string>

  <key>RunAtLoad</key>       <true/>
  <key>KeepAlive</key>       <true/>

  <key>StandardOutPath</key> <string>/Users/jaywest/Library/Logs/claudeclaw/dashboard.out</string>
  <key>StandardErrorPath</key><string>/Users/jaywest/Library/Logs/claudeclaw/dashboard.err</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
