# infra/launchd/com.claudeclaw.decay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>           <string>com.claudeclaw.decay</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/jaywest/claudeclaw/.venv/bin/python</string>
    <string>-m</string>
    <string>memory.decay</string>
  </array>

  <key>WorkingDirectory</key><string>/Users/jaywest/claudeclaw</string>

  <!-- Sundays at 03:30 local — runs after the 03:15 backup completes. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key> <integer>0</integer>
    <key>Hour</key>    <integer>3</integer>
    <key>Minute</key>  <integer>30</integer>
  </dict>

  <key>RunAtLoad</key>       <false/>
  <key>StandardOutPath</key> <string>/Users/jaywest/Library/Logs/claudeclaw/decay.out</string>
  <key>StandardErrorPath</key><string>/Users/jaywest/Library/Logs/claudeclaw/decay.err</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>DB_PATH</key>
    <string>/Users/jaywest/claudeclaw/store/claudeclaw.db</string>
  </dict>
</dict>
</plist>
